# deploy/supervisor.conf
; Supervisor program config for the Used Vehicle Parts Search app.
;
; Install this file at /etc/supervisor/conf.d/uvps.conf (requires sudo)
; then: sudo supervisorctl reread && sudo supervisorctl update
;
; Logs land in /var/log/supervisor/uvps-{stdout,stderr}.log
;
; Replace these placeholders before installing:
;   {APP_USER}     -> the Linux user that owns the app dir
;                     (on Cloudways, your application's master user, e.g. master)
;   {APP_DIR}      -> absolute path to the project root
;                     (e.g. /home/master/applications/abcdef/private_html/uvps)
;   {VENV_DIR}     -> absolute path to the virtualenv
;                     (e.g. {APP_DIR}/.venv)
;   {PORT}         -> bind port for uvicorn (e.g. 8000)

[program:uvps]
command={VENV_DIR}/bin/uvicorn app.main:app --host 127.0.0.1 --port {PORT} --workers 2
directory={APP_DIR}
user={APP_USER}
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=/var/log/supervisor/uvps-stdout.log
stderr_logfile=/var/log/supervisor/uvps-stderr.log
environment=PYTHONUNBUFFERED="1"
