# deploy/supervisor.conf
; Supervisor program config for the Used Vehicle Parts Search app.
;
; deploy/deploy.sh installs this automatically. If you need to install manually:
;   sudo cp deploy/supervisor.conf /etc/supervisor/conf.d/uvps.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Replace these placeholders before installing manually:
;   {APP_USER}  -> Linux user that owns the app (default: uvps)
;   {APP_DIR}   -> absolute path to the project root (e.g. /home/uvps/uvps)
;   {VENV_DIR}  -> absolute path to the virtualenv  (e.g. /home/uvps/uvps/.venv)
;   {PORT}      -> uvicorn bind port                 (e.g. 8000)

[program:uvps]
command={VENV_DIR}/bin/uvicorn app.main:app --host 127.0.0.1 --port {PORT} --workers 2
directory={APP_DIR}
user={APP_USER}
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=/var/log/supervisor/uvps-stdout.log
stderr_logfile=/var/log/supervisor/uvps-stderr.log
environment=PYTHONUNBUFFERED="1"
